# vault-indexer/crontab
0 2 * * * python /app/indexer.py >> /state/indexer.log 2>&1
